# edge-server/crontab
* * * * * curl http://localhost:80/process-cache